# hack/docker/s6/gin-music/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

export USER=ginmusic
exec gosu $USER /app/ginmusci/gin-music serv --config /data/ginmusic/conf/app.toml
